# Documents/PowerShell/profile.d/50-navi.ps1
if (Get-Command navi -ErrorAction SilentlyContinue) {
    navi widget powershell | Out-String | Invoke-Expression
}
