# Documents/PowerShell/profile.d/90-atuin.ps1
if (Get-Command atuin -ErrorAction SilentlyContinue) {
    atuin init powershell | Out-String | Invoke-Expression
}
